# deploy/splitkostya.service
[Unit]
Description=SplitKostya Angular Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/splitKostya/current
ExecStart=/usr/bin/node dist/splitKostya/server/server.mjs
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
